# src/temperature-sensor/config/smart-pc-temperature-sensor@.service.in
[Unit]
Description=Smart PC Control - Temperature Daemon
After=syslog.target network.target

[Service]
User=%i
Type=forking
EnvironmentFile=-/@CMAKE_INSTALL_SYSCONFIG_PATH@/smart-pc-temperature-sensor
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/smart-pc-temperature-sensor -d -b $MQTT_BROKER_URL -i $CLIENT_ID $DAEMON_OPTIONS
PrivateTmp=false
Restart=on-failure

[Install]
WantedBy=multi-user.target