# src/power-control/config/smart-pc-control-power@.service.in
[Unit]
Description=Smart PC Control - Power Daemon
After=syslog.target network.target

[Service]
User=%i
Type=forking
EnvironmentFile=-/@CMAKE_INSTALL_SYSCONFIG_PATH@/smart-pc-control-power
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/smart-pc-control-power -d -b $MQTT_BROKER_URL -i $CLIENT_ID $DAEMON_OPTIONS
PrivateTmp=false

[Install]
WantedBy=multi-user.target